# scripts/systemd/lumora-backup.timer
[Unit]
Description=Run the Lumora data backup daily

[Timer]
# Daily at 03:30. Persistent so a missed run (VPS off) fires on next boot.
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
